# systemd/ground_loop.user.service
[Unit]
Description=UnseenUniversity Ground Loop — plugin-host process supervisor (user unit)
Documentation=https://github.com/akienm/UnseenUniversity
# User-manager variant of ground_loop.service. Runs without sudo and starts
# at boot because the user has lingering enabled (loginctl enable-linger akien).
# Install:
#   cp systemd/ground_loop.user.service ~/.config/systemd/user/ground_loop.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now ground_loop.service
After=network.target

[Service]
Type=simple

# Repo root — adjust if installed elsewhere
WorkingDirectory=/home/akien/dev/src/UnseenUniversity

# Auto-deploy plugin descriptors from the repo into the live scan dir on every
# start. Without this, the descriptor only lives in config/ and the daemon
# (which scans ~/.unseen_university/ground_loop/) never sees it — the exact
# failure that took the web server down. Makes the supervisor self-deploying.
# Fail-OPEN: a copy failure must never block the supervisor from starting (it
# matches Ground Loop's fail-open philosophy and the "never down" intent — a
# stale-but-present descriptor beats no supervisor). The `|| true` guarantees
# ExecStartPre always succeeds; worst case the daemon scans whatever descriptors
# are already deposited.
ExecStartPre=-/bin/sh -c 'mkdir -p %h/.unseen_university/ground_loop && cp /home/akien/dev/src/UnseenUniversity/config/ground_loop/*.yaml %h/.unseen_university/ground_loop/ 2>/dev/null || true'

# Python from the repo's venv
ExecStart=/home/akien/dev/src/UnseenUniversity/.venv/bin/python3 \
          -m devices.ground_loop.daemon \
          --poll 15 \
          --log-level INFO

# Restart on failure (systemd keeps it alive; Ground Loop restarts plugins internally)
Restart=on-failure
RestartSec=10

# Environment
Environment=UU_ROOT=/home/akien/dev/src/UnseenUniversity
EnvironmentFile=-%h/.unseen_university/uu.env

# Logs go to the user journal — view with: journalctl --user -u ground_loop -f
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
